# org.kevoree.watchdog/src/deb/control/postrm
#!/bin/sh
sudo update-rc.d kevoree remove

if [ -n "$(getent passwd kevoree)" ];
    then userdel kevoree;
fi

if [ -n "$(getent group kevoree)" ];
    then groupdel kevoree
fi
