# org.kevoree.watchdog/src/deb/control/postinst
#!/bin/sh
chmod +x /etc/init.d/kevoree
echo "Registering kevoree service"
update-rc.d kevoree defaults
