# org.kevoree.watchdog/src/deb/control/postrm
#!/bin/sh
sudo update-rc.d remove kevoree
