# org.kevoree.watchdog/src/deb/control/prerm
#!/bin/sh
echo "Stopping kevoree service"
service kevoree stop
echo "Removing kevoree service"
update-rc.d -f kevoree remove
